# etc/crontab
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

0 * * * *   root	( /sbin/apcaccess | /var/lib/app-apc-reporter/status ) > /dev/null

